# deploy/hcmos.service
[Unit]
Description=HCMOS/ESS API + frontend (UAT)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=hcmos
WorkingDirectory=/opt/hcmos
EnvironmentFile=/etc/hcmos/hcmos.env
# NODE_ENV=production (from the env file) disables the test-only fault seams.
ExecStart=/usr/bin/node src/server.js
Restart=on-failure
RestartSec=3
# Hardening: no new privileges, private tmp, read-only system dirs.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/hcmos
# KIOSK photo-on-punch (Kira 2026-07-14): punch photos are binary at rest
# OUTSIDE the deploy tree. StateDirectory makes systemd create
# /var/lib/hcmos/punch-photos owned by the service user and writable inside
# the ProtectSystem=strict sandbox — without this the photo write fails and
# every kiosk punch records flagged no-photo (it degrades, never blocks).
StateDirectory=hcmos/punch-photos
StateDirectoryMode=0700
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
